# deploy/encoder.service
[Unit]
Description=ZuidWest FM Encoder
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/encoder -config /etc/encoder/config.json
Restart=always
RestartSec=5
User=root

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/etc/encoder

[Install]
WantedBy=multi-user.target
