# deploy/encoder.service
[Unit]
Description=ZuidWest FM Encoder
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/encoder -config /etc/encoder/config.json
Restart=always
RestartSec=5

# Run as dedicated service user
User=encoder
Group=encoder

# Systemd-managed directories
LogsDirectory=encoder
LogsDirectoryMode=0750

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/etc/encoder

[Install]
WantedBy=multi-user.target
